# config/com.claude.monitor.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service identification -->
    <key>Label</key>
    <string>com.claude.monitor.daemon</string>
    
    <!-- Program execution -->
    <key>Program</key>
    <string>__PYTHON_PATH__</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON_PATH__</string>
        <string>__DAEMON_PATH__/run_daemon.py</string>
        <string>--start-day</string>
        <string>__BILLING_START_DAY__</string>
        <string>--interval</string>
        <string>10</string>
        <string>--time-alert</string>
        <string>30</string>
        <string>--inactivity-alert</string>
        <string>10</string>
        <string>--sessions</string>
        <string>50</string>
        <string>--timezone</string>
        <string>Europe/Warsaw</string>
    </array>
    
    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>__DAEMON_PATH__</string>
    
    <!-- Automatic startup configuration -->
    <key>RunAtLoad</key>
    <true/>
    
    <!-- Keep daemon running -->
    <key>KeepAlive</key>
    <true/>
    
    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__USER_HOME__</string>
        <key>PATH</key>
        <string>__FULL_PATH__</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>USER</key>
        <string>__USER_NAME__</string>
    </dict>
    
    <!-- Logging configuration -->
    <key>StandardOutPath</key>
    <string>__USER_HOME__/.config/claude-monitor/daemon.log</string>
    
    <key>StandardErrorPath</key>
    <string>__USER_HOME__/.config/claude-monitor/daemon.error.log</string>
    
    <!-- Resource limits -->
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>1024</integer>
        <key>NumberOfProcesses</key>
        <integer>128</integer>
    </dict>
    
    <!-- Hard resource limits to prevent runaway processes -->
    <key>HardResourceLimits</key>
    <dict>
        <key>NumberOfProcesses</key>
        <integer>256</integer>
    </dict>
    
    <!-- Throttling configuration -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <!-- Session type -->
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>
    
    <!-- Process type -->
    <key>ProcessType</key>
    <string>Standard</string>
    
    <!-- Exit timeout -->
    <key>ExitTimeOut</key>
    <integer>30</integer>
    
    <!-- Nice level -->
    <key>Nice</key>
    <integer>1</integer>
    
    <!-- Low priority background process -->
    <key>LowPriorityIO</key>
    <true/>
    
</dict>
</plist>